# supervisor-vsftpd.conf
[program:vsftpd]
command=/usr/sbin/run-vsftpd.sh
